# integration/8.2-integration/docker-integration.conf
[program:php-fpm]
command=php-fpm -R
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
startretries=3
autorestart=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0

[program:cron]
command=cron -f
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
startretries=3
autorestart=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
